# build/Lede_phicomm_n1/files/etc/crontabs/root
0 */3 * * * /root/DenyPwdHack.sh
*/20 * * * * /root/ddns.sh
